# src/service.in
#!/bin/bash
#
# applefand         Start/Stop the apple fan daemon
#
# chkconfig: 2345 50 50
# description: applefand monitors CPU temperature and adjust fan speed
# accordingly

RETVAL=0
prog="applefand"
prefix=@PREFIX@
exec="$prefix/bin/$prog"
lockfile="/var/lock/subsys/$prog"
PREFIX=$PREFIX:$prefix

# Source function library
. /etc/rc.d/init.d/functions

start() {
    [ -x $exec ] || exit 5
    echo -n $"Starting $prog: "
    daemon $exec -d
    retval=$?
    echo
    [ $retval -eq 0 ] && touch $lockfile
}

stop() {
    echo -n "Stopping $prog: "
    if [ -n "$(ps aux | grep $prog | grep -v grep | awk '{ print $2 }' | head -n 1)" ] ; then
        killproc $exec
        RETVAL=3
    else
        failure $"Stopping $prog"
    fi
    retval=$?
    echo
    [ $retval -eq 0 ] && rm -f $lockfile
}

restart() {
    stop
    start
}

rh_status() {
    pid="$(ps aux | grep $exec | grep -v grep | awk '{ print $2 }' | head -n 1)"
    if [ -n "$pid" ] ; then
        echo $"$prog (pid $pid) is running..."
        retval=0
    else
        echo $"$prog is stopped"
        retval=3
    fi
    return $retval
}

rh_status_q() {
    rh_status >/dev/null 2>&1
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        restart
        ;;
    status)
        rh_status
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart|status}"
        exit 1
esac

exit $?

